# owl/s500/burn/card_burn_initramfs/etc/inittab
# /etc/inittab
::sysinit:/etc/init.d/rcS
/dev/console::respawn:-/bin/sh
/dev/ttyS2::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r

