# owl/s500/burn/initramfs/etc/inittab
# /etc/inittab
::sysinit:/etc/init.d/rcS
/dev/console::askfirst:-/bin/sh
#ttyS1::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r

